# master1/dbus-org.freedesktop.nm-dispatcher.service
[Unit]
Description=Network Manager Script Dispatcher 